# systemd/user/fehbg.timer
[Unit]
Description=Change wallpapers every 10 minutes

[Timer]
OnUnitActiveSec=10min
