# systemd/user/mbsync.timer
[Unit]
Description=Sync email every hour

[Timer]
OnUnitActiveSec=1hour

[Install]
WantedBy=timers.target